# HackChat/API-Server/hackchat.api-com-1.service
[Unit]
Description=HackChat COM API Service

[Service]
Type=simple
Restart=always
RestartSec=15
ExecStart=bash /usr/bin/hackchat-com-api.sh
